# openwrt/proxy/init.d/ssh-http
#!/bin/sh /etc/rc.common  
# Copyright (C) 2007 OpenWrt.org  
START=99

start() {
    polipo -c '/etc/polipo/config'
}

stop() {  
    killall polipo  
}