# server/router/files/files/etc/uci-defaults/99-custom
uci -q batch << EOI
set system.@system[0].zonename='Asia/Shanghai'
set system.@system[0].timezone='CST-8'
commit system

#set network.lan.ipaddr='192.168.2.15'
#set network.lan.gateway='192.168.2.2'
#set network.lan.dns='114.114.114.114'
#commit network

# open wifi
#set wireless.radio0.disabled='0'
#commit wireless

EOI
